# scripts/caso.cron
10 * * * * root /usr/local/bin/run-caso.sh

