# scripts/ssmsend.cron
01 01 * * * root /usr/local/bin/run-ssmsend.sh

